# templates/default/h2o.service.erb
[Unit]
Description=h2o - the optimized HTTP/1.x, HTTP/2 server
Documentation=http://h2o.org/en/docs/
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=simple
PIDFile=<%= node['h2o']['pid'] %>
ExecStart=<%= node['h2o']['bin'] %>/h2o -m master -c <%= node['h2o']['dir'] %>/h2o.conf
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
PrivateTmp=<%= node['h2o']['systemd_private_tmp'] ? 'true' : 'false' %>
LimitNOFILE=infinity

[Install]
WantedBy=multi-user.target
